# deploy/vaultr-ebay-webhook.service
[Unit]
Description=Vaultr eBay webhook listener
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/Documents/GitHub/vaultr
ExecStart=/usr/bin/node /home/pi/Documents/GitHub/vaultr/dist/ebay-webhook.js
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
