# deploy/vaultr-ebay-webhook.service
[Unit]
Description=Vaultr eBay webhook listener
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/Documents/GitHub/vaultr
ExecStartPre=/usr/bin/mkdir -p /home/pi/Documents/GitHub/vaultr/data/logs
ExecStart=/usr/bin/node /home/pi/Documents/GitHub/vaultr/dist/ebay-webhook.js
Restart=always
RestartSec=3
MemoryAccounting=yes
MemoryHigh=150M
MemoryMax=250M
StandardOutput=append:/home/pi/Documents/GitHub/vaultr/data/logs/ebay-webhook.log
StandardError=append:/home/pi/Documents/GitHub/vaultr/data/logs/ebay-webhook.log

[Install]
WantedBy=multi-user.target
